# kafka.conf
[program:kafka]
autostart=true
autorestart=true
