# kafka.conf
[program:kafka]
command=/opt/kafka_2.10-0.8.2.1/bin/kafka-server-start.sh /opt/kafka_2.10-0.8.2.1/config/server.properties
autostart=true
autorestart=true
